# etc/profile.d/dwm-conf.sh
export LD_LIBRARY_PATH="$HOME/.config/phyos/pdwm"
